# playbook/roles/middleware/neo4j/files/neo4j.conf
description "neo4j"
author  "70_10 <example@mail.com>"

start on startup
stop on shutdown

chdir /opt/neo4j
exec ./bin/neo4j start
respawn
